# node/openclaw-node.service
[Unit]
Description=OpenClaw Node Host (Docker Compose)
Requires=docker.service
After=docker.service

[Service]
Type=simple
User=openclaw
WorkingDirectory=/home/openclaw
ExecStartPre=/usr/local/bin/docker-compose pull
ExecStart=/usr/local/bin/docker-compose up
ExecStop=/usr/local/bin/docker-compose down
Restart=always

[Install]
WantedBy=multi-user.target
